# docker/supervisor/services/friendship.ini
[program:friendship]
command=php bin/console messenger:consume -vv friendship --limit=1000
redirect_stderr=true
stderr_syslog=true
numprocs=1
directory=/var/www
autostart=true
autorestart=true
startretries=1000
user=www-data
stopsignal=KILL
